# Ooziee-Job-Reprocessing-Utility/deploy/systemd/oozie-reprocess-worker.service
[Unit]
Description=Oozie Reprocessing Manager Worker
After=network.target redis.service mysqld.service mariadb.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/oozie-reprocessing-manager/worker
EnvironmentFile=/opt/oozie-reprocessing-manager/backend/.env
ExecStart=/opt/oozie-reprocessing-manager/worker/.venv/bin/python runner.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
