# Ooziee-Job-Reprocessing-Utility/deploy/systemd/oozie-reprocess-worker.service
[Unit]
Description=Oozie Reprocessing Manager Worker
After=network-online.target redis.service mysqld.service mariadb.service
Wants=network-online.target

[Service]
Type=simple
User=ooziemgr
Group=ooziemgr
WorkingDirectory=/opt/oozie-reprocessing-manager/worker
EnvironmentFile=/etc/oozie-reprocessing/oozie-reprocess.env
Environment=PYTHONPATH=/opt/oozie-reprocessing-manager/backend
ExecStart=/opt/oozie-reprocessing-manager/worker/.venv/bin/python runner.py
TimeoutStartSec=30
TimeoutStopSec=30
Restart=always
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
